# scripts/fa.service
[Unit]
# systemd USER unit (installed to ~/.config/systemd/user/fa.service).
# NOTE: a user unit cannot Requires=/After= a *system* unit such as
# docker.service — systemd resolves dependencies within the user manager only.
# Boot ordering against Docker is handled by the container's own
# `restart: unless-stopped` policy plus Docker's `live-restore`. We only order
# softly against networking so `up -d` does not race a cold network.
Description=First-Agent 24/7 container
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
OOMScoreAdjust=-500
Restart=on-failure
RestartSec=10
WorkingDirectory=/srv/first-agent/repo/First-Agent-dev
ExecStart=/usr/bin/docker compose -f docker-compose.fa.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.fa.yml down
TimeoutStartSec=0

[Install]
WantedBy=default.target
